# etc/systemd/system/slipper.timer
[Unit]
Description=Slipper Timer

[Timer]
OnBootSec=5min
OnCalendar=*:0/1
Unit=slipper.target

[Install]
WantedBy=basic.target
