# etc/systemd/system/slipper.target
[Unit]
Description=Slipper Timer Target
StopWhenUnneeded=yes
